# recipes-bumperbot/bumperbot-systemd/bumperbot-systemd/bumperbot.target
[Unit]
Description=Bumperbot ROS 2 Stack Target
Wants=network-online.target
After=network-online.target

[Install]
WantedBy=multi-user.target
