# recipes-bumperbot/bumperbot-systemd/bumperbot-systemd/bumperbot.target
[Unit]
Description=Bumperbot ROS 2 Stack Target
Requires=multi-user.target
After=network.target

[Install]
WantedBy=multi-user.target
