# recipes-core/udev/udev-extraconf/99-wifi-unblock.rules
SUBSYSTEM=="net", ACTION=="add", KERNEL=="wlan0", RUN+="/usr/sbin/rfkill unblock wifi"
